# scripts/systemd/social-collect.service
[Unit]
Description=Daily social-update collection (enqueues a run for the local poller)
After=network-online.target

[Service]
Type=oneshot
# Collectors now run through the queue: the daily timer just enqueues a run and the
# social-collect-poll unit claims it and does the actual collection. Single path.
ExecStart=/bin/bash -lc '%h/.local/bin/social-collect-poll.sh --enqueue daily'
SyslogIdentifier=social-collect
