# scripts/systemd/social-collect.service
[Unit]
Description=Social-update collector (watchdog-style daily run)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -lc '%h/.local/bin/social-collect-watchdog.sh'
SyslogIdentifier=social-collect
