# scripts/systemd/social-collect-poll.service
[Unit]
Description=Social-update collection poller (claims UI/timer-triggered runs)
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash -lc '%h/.local/bin/social-collect-poll.sh'
SyslogIdentifier=social-collect-poll
# The watchdog closes the CDP Chrome it launches at the end of each run, so no
# KillMode override is needed — default control-group cleanup is a fine backstop.
