# scripts/systemd/vault-daily.timer
[Unit]
Description=Run the vault-keeper daily writer each evening (after the 18:00 collector)

[Timer]
OnCalendar=*-*-* 19:00:00
# run after a missed window (WSL was down at 19:00) instead of skipping the day
Persistent=true

[Install]
WantedBy=timers.target
